# meta-renesas-ivi/recipes-core/busybox/busybox-1.21.1/ntpd.service
[Unit]
Description=Network Time Protocol client
Before=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/ntpd-set.sh
RemainAfterExit=true
